# systemd/eink-display.service
[Unit]
Description=E-Ink Pet Clock Display Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%u
Group=%u
WorkingDirectory=%h/einkpetclock
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 %h/einkpetclock/core/display_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
Nice=-5
CPUQuota=80%

[Install]
WantedBy=multi-user.target
