# systemd/eink-api.service
[Unit]
Description=E-Ink Pet Clock API Service
After=network.target eink-display.service
Wants=network-online.target

[Service]
Type=simple
User=%u
Group=%u
WorkingDirectory=%h/einkpetclock/web
Environment="PYTHONUNBUFFERED=1"
ExecStart=%h/einkpetclock/venv/bin/python %h/einkpetclock/web/api.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
CPUQuota=50%

[Install]
WantedBy=multi-user.target
